# peripheral-manager/xyz.openbmc_project.peripheral.manager.service.in
[Unit]
Description= NIC management
Wants=obmc-mapper.target
After=obmc-mapper.target

[Service]
ExecStart=@bindir@/peripheral-manager
Type=dbus
BusName=xyz.openbmc_project.peripheral.manager
SyslogIdentifier=peripheral-manager

[Install]
WantedBy=multi-user.target
